# OneDrive/Escritorio/-=L.R.GONZALEZ=-/packages/open-dev/claudio-os-blueprint/live-build/config/includes.chroot/etc/systemd/system/claudio-dashboard.service
[Unit]
Description=Claudio Mission Control Dashboard
After=claudio-guardian.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/claudio/dashboard/dashboard.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

